# DEBIAN/postrm
#!/bin/sh -e

systemctl stop getty@tty12.service
systemctl disable getty@tty12.service
systemctl daemon-reload

deluser --quiet --system pkgsync

exit 0
