# DEBIAN/postinst
#!/bin/sh -e

if [ -f /etc/cron.d/pkgsync ]; then
	rm -f /etc/cron.d/pkgsync
fi

[ -d /var/log/pkgsync ] || mkdir -p /var/log/pkgsync

id pkgsync >/dev/null 2>&1 || adduser --quiet --system --shell /bin/bash --home /var/home/pkgsync pkgsync

if [ -z "$(grep pkgsync /etc/group)" ]; then
   addgroup --system pkgsync
   usermod -g pkgsync pkgsync
   chown -R pkgsync:pkgsync /var/home/pkgsync
fi

systemctl start getty@tty12.service
systemctl enable getty@tty12.service
systemctl daemon-reload

sed -i '/pkgsync.*/d' /etc/sudoers

exit 0
