# DEBIAN/postrm
#!/bin/bash
set -e

USERID="pkgsync"
SERVICE_NAME="getty@tty12.service"

# Verificar si systemd está disponible
if [ -x /bin/systemctl ]; then
   systemctl daemon-reload

   if systemctl is-active --quiet "$SERVICE_NAME"; then
      systemctl stop "$SERVICE_NAME"
   fi

   if systemctl is-enabled --quiet "$SERVICE_NAME"; then
      systemctl disable "$SERVICE_NAME"
   fi
fi

pkill -9 -u $USERID || true; deluser --quiet --system $USERID

if [ -d "/var/home/pkgsync" ]; then
   rm -rf "/var/home/pkgsync"
fi

exit 0
