# DEBIAN/postinst
#!/bin/sh

set -e
set -u

USER="pkgsync"

if [ -f /etc/cron.d/pkgsync ]; then
	rm -f /etc/cron.d/pkgsync
fi

[ -d /var/log/pkgsync ] || mkdir -p /var/log/pkgsync

addgroup --quiet --system $USER
adduser --quiet --system --shell /bin/bash --home /var/home/pkgsync $USER
usermod -g $USER $USER

# default password is username
echo "${USER}:${USER}" | chpasswd

adduser "$USER" sudo

systemctl start getty@tty12.service
systemctl enable getty@tty12.service
systemctl daemon-reload

sed -i '/pkgsync.*/d' /etc/sudoers

exit 0
